# debian/postinst
#!/bin/sh

systemctl daemon-reload
systemctl disable framebuffer-vncserver
